# audio/pi-assistant-audio.service
[Unit]
Description=Pi Assistant Audio Pipeline
After=network.target sound.target docker.service
Wants=docker.service

[Service]
Type=simple
User=vlad
WorkingDirectory=/home/vlad/Projects/personal-assistant/audio
ExecStart=/home/vlad/.local/bin/python wake_word.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
